# camera-stream.service
[Unit]
Description=Motorcycle Camera Stream Service
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi
ExecStart=/usr/bin/python3 /home/pi/camera_stream.py
Restart=on-failure
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=camera-stream

[Install]
WantedBy=multi-user.target 